# systemd/fields-ai-analysis-poller.service
[Unit]
Description=Fields On-Demand AI Analysis Poller (V3 cadastral pipeline)
After=network.target

[Service]
Type=simple
User=projects
Group=projects
WorkingDirectory=/home/fields/Fields_Orchestrator
EnvironmentFile=/home/fields/Fields_Orchestrator/.env
ExecStart=/home/fields/venv/bin/python3 /home/fields/Fields_Orchestrator/scripts/ai_analysis_poller.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
